# scripts/supervisor-gohop.ini
[program:gohop]
command=/usr/local/bin/gohop /etc/gohop/server.ini
autostart=true
autorestart=true
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gohop.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/gohop.error
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
